# cron.daily/emaint_check
#!/bin/bash

# Copyright Douglas J Hunley 2013

# Call emaint and have it fix everything it can
/usr/sbin/emaint --fix all &>/dev/null
