# cron.daily/emerge_world
#!/bin/bash

if [ ! -s /tmp/emerge.out ] ; then
    exit 0
fi

emerge \
  --quiet \
  --update \
  --newuse \
  --deep \
  --with-bdeps y \
  --keep-going \
  @system @world &>/dev/null
